# docker/etc/supervisor/conf.d/php7.3-fpm.conf
[program:php7.3-fpm]
command=/usr/sbin/php-fpm7.3 -F
priority=9
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
